# sv/window_title/run
#!/bin/sh

exec 2>&1

echo "=== Starting sv/window_title: $(date)"

set -x
cd /apps/my_desktop
exec bin/my_desktop window_title