# sv/publish_pulseaudio/run
#!/bin/sh

exec 2>&1

echo "=== Starting sv/publish_pulseaudio: $(date)"

set -x
cd /apps/my_desktop
exec bin/my_desktop publish_pulseaudio
